# cron/app.cron
*/5 * * * * . $HOME/.profile; mkdir -p "$DATA_DIR/log/cron/" && "$MAIN_DIR"/run backup >> "$DATA_DIR/log/cron/backup.`date '+\%F'`.log" 2>&1
*/1 * * * * . $HOME/.profile; mkdir -p "$DATA_DIR/log/cron/" && "$MAIN_DIR"/run "sync:verify" >> "$DATA_DIR/log/cron/sync.`date '+\%F'`.log" 2>&1
